# deployment/jadzia.service
[Unit]
Description=Jadzia V4 - WordPress Management Agent
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=jadzia
Group=jadzia
WorkingDirectory=/opt/jadzia

# Environment
Environment="PATH=/opt/jadzia/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/opt/jadzia/.env

# Start command (production — no reload; use UVICORN_RELOAD=1 only for local dev via main.py)
ExecStart=/opt/jadzia/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000

# Restart policy
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitInterval=300

# Logging
StandardOutput=append:/opt/jadzia/logs/jadzia.log
StandardError=append:/opt/jadzia/logs/jadzia-error.log

# Resource limits
LimitNOFILE=65536
MemoryMax=2G
CPUQuota=200%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
